# arch/install.sh
#!/bin/bash

sudo pacman -Sy networkmanager network-manager-applet networkmanager-openvpn nvim dbus i3-gaps tmux dmenu compton mupdf ranger mopidy ncmpcpp mpc mpv mplayer dunst cronie feh scrot ffmpeg w3m xcape youtube-dl newsboat wget curl arandr firefox chromium gimp inkscape xclip xorg gcc make fontconfig pkg-config lxappearance papirus-icon-theme alsa-utils imagemagick ntfs-3g unzip gnome-themes-standard gtk-engine-murrine telegram-desktop qt5ct r udisks2 xdotool xssstate sxiv

sudo systemctl enable NetworkManager
